# add-ons/extensions/usb-exec/10-usbplugin.rules.template
KERNEL=="sg*[0-9]", ACTION=="add", RUN+="{{EXECUTABLE}}"
